# supervisor.conf
[program:mimi]
command = /home/ubuntu/mimi/mimi_start.sh

; General settings
autostart = true
autorestart = true
