# supervisor.conf
[program:mimi]
command = /home/ubuntu/mimi/mimi_start.sh

; General settings
autostart = true
autorestart = true

; Kill all in group
killasgroup = true

; Low priority
priority = 1000
